# luci-theme-vdesign/root/etc/uci-defaults/30_luci-theme-design
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.vDesign=/luci-static/vdesign
		set luci.main.mediaurlbase=/luci-static/vdesign
		commit luci
	EOF
fi

exit 0